# board/myrpi4b/rootfs_overlay/etc/profile.d/variables.sh
#!/bin/bash

SETTINGS='/data/settings.txt'
MY_VAR='test_my_variable'
